# backends/upstart/live-config.upstart
## Debian Live - System Configuration Scripts
##
## live-config contains the scripts that configure a Debian Live system during
## the boot process (late userspace).
##
## This is the upstart job for live-config.


description	"live-config"
author		"Debian Live Project <debian-live@lists.debian.org>"

start on filesystem
stop on runlevel [016]

console output

script
	if [ -e /lib/live/config.sh ]
	then
		/lib/live/config.sh
	fi
end script
